# web/drush-archive/code/cronjob_backup/crontab
0 0 * * * /var/www/html/cronjob_backup/backup.sh
